# sample/crontab
* * * * *   /usr/bin/execute /sample.sh

